# Documents/quickrevive.service
[Unit]
Description=You need a little a revive

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sysd-check.sh
User=root
